# com.hive.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.hive.orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npx</string>
        <string>ts-node</string>
        <string>/Users/carlos.miranda/Documents/Github/hive/orchestrator.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/carlos.miranda/Documents/Github/hive</string>

    <!-- Run at midnight every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>DATABASE_URL</key>
        <string>postgresql://neondb_owner:npg_sT6mLbxz5BAp@ep-frosty-meadow-ag9xf51n-pooler.c-2.eu-central-1.aws.neon.tech/neondb?sslmode=require</string>
        <key>ENCRYPTION_KEY</key>
        <string>741e1e9fcc1795da6daaf2ee9d5b212d8e52540637f6f92c7a8da521d400e92e</string>
        <key>HIVE_API_URL</key>
        <string>https://hive-phi.vercel.app/api</string>
        <key>NEXT_PUBLIC_URL</key>
        <string>https://hive-phi.vercel.app</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/carlos.miranda/Documents/Github/hive/logs/orchestrator.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/carlos.miranda/Documents/Github/hive/logs/orchestrator.stderr.log</string>

    <!-- No KeepAlive — run once at midnight, don't retry on failure.
         Errors are logged and the self-healing cycle handles them next run. -->
</dict>
</plist>
